# config/systemd/music-app.target
[Unit]
Description=Music App
Requires=music-app.target systemd-networkd.service
After=music-app.target systemd-networkd.service
AllowIsolate=yes
