# config/systemd/puma.target
[Unit]
Description=Puma
Requires=puma.target systemd-networkd.service
After=puma.target systemd-networkd.service
AllowIsolate=yes
