# config/systemd/sidekiq.target
[Unit]
Description=Sidekiq
Requires=sidekiq.target systemd-networkd.service
After=sidekiq.target systemd-networkd.service
AllowIsolate=yes
